# config/service_asker_bot.conf
[program:service_asker_bot]
command=/www/service_bot_venv/bin/python main.py
directory=/www/service_bot/
user=www-data
autorestart=true